# launch-agents/net.graphify.labirynt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>net.graphify.labirynt</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/michalmaciejewski/.claude/scripts/graphify-labirynt-weekly.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>17</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/michalmaciejewski/.claude/scripts/logs/graphify-launchd.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/michalmaciejewski/.claude/scripts/logs/graphify-launchd.err</string>

    <key>WorkingDirectory</key>
    <string>/Users/michalmaciejewski/Desktop/Labirynt</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/michalmaciejewski</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
